# ustc-checkin.timer
[Unit]
Description=USTC Daily Checkin Timer

[Timer]
OnCalendar=*-*-* 10:00
RandomizedDelaySec=10h
Persistent=true

[Install]
WantedBy=timers.target
